# deploy/led-panel-emergency-restart.service.example
# Standalone MQTT-triggered watchdog for led-panel-controller.service. Deliberately a
# separate service/process from the controller itself -- if the controller is what's
# hung, it can't act on its own MQTT commands (its single dispatch loop is what's
# stuck), so the thing that force-restarts it has to live outside it.
#
# Install: See README.md in this directory for instructions on how to install and enable the service.
#
# Then in Home Assistant: button.led_display_emergency_restart appears automatically
# via MQTT Discovery, same device as the rest of the LED Display Controller entities.

[Unit]
Description=LED Panel Controller Emergency Restart Watchdog
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/led-panel-controller/deploy/emergency-restart-watchdog.sh
EnvironmentFile=-/opt/led-panel-controller/deploy/led-panel-emergency-restart.env
Restart=always
RestartSec=5
# Needs root to run `systemctl restart led-panel-controller.service`.
User=root

[Install]
WantedBy=multi-user.target
